# code/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:databasesetup]
priority=1
command=SetupDatabase
autorestart=false

#[program:recommender]
#priority=10
#command=RunRecommenderSystem
#autorestart=true